# kalshi_daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kalshi.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/anaconda3/bin/python3</string>
    <string>/Users/ankhbayarbatkhurel/kalshi-calibration/src/run_daily.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ankhbayarbatkhurel/kalshi-calibration</string>

  <!-- Run at 10am, 1pm, and 3pm EDT (14:00, 17:00, 19:00 UTC) -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>   <integer>14</integer>
      <key>Minute</key> <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>   <integer>17</integer>
      <key>Minute</key> <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>   <integer>19</integer>
      <key>Minute</key> <integer>0</integer>
    </dict>
  </array>

  <!-- Delay launch until DNS/network is available (macOS 11+) -->
  <key>WaitForNetworkConnectivity</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/ankhbayarbatkhurel/kalshi-calibration/data/trading.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ankhbayarbatkhurel/kalshi-calibration/data/trading_error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
